# src/modules/usage-statistics/filesystem/root/etc/systemd/system/usage-statistics.service
[Unit]
Description=Send anonymous statistics of a distro
After=multi-user.target
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/boot_report
[Install]
WantedBy=multi-user.target
